# deploy/dashboard.service
[Unit]
Description=Amber Terminal Dashboard for CRT Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/dashboard
Restart=always
RestartSec=5
Environment=TERM=linux
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
